# telegram-server-monitor/security-monitor.service
[Unit]
Description=Telegram Server Monitor
After=network.target

[Service]
ExecStart=/usr/bin/bash /home/mark/telegram-server-monitor/monitor.sh
